# rc.local
#!/bin/sh -e

_IP=$(hostname -I) || true



exit 0
